# app/scripts/opsdeck.service
# systemd USER service (opsional) — alternatif autostart untuk Linux.
# Cocok bila ingin app di-restart otomatis kalau crash.
#
# Pasang:
#   mkdir -p ~/.config/systemd/user
#   cp opsdeck.service ~/.config/systemd/user/
#   (edit ExecStart ke path AppImage/binary hasil install)
#   systemctl --user daemon-reload
#   systemctl --user enable --now opsdeck.service
#   loginctl enable-linger $USER   # biar jalan tanpa harus login GUI (kalau perlu)

[Unit]
Description=Opsdeck
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/opt/Opsdeck.AppImage
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
